# confs/celery.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid ubuntu
setgid ubuntu

script
  /bin/bash /home/ubuntu/apps/{{ REPOSITORY_NAME }}/run_on_server.sh celery worker --app={{ PROJECT_FILENAME }}.celery.app -l info -f {{ SERVER_LOG_PATH }}/celery.log
end script